# udev/95-monitor-hot-plug.rules
SUBSYSTEM=="drm", RUN+="su LOGICIFY\\corvis -c /usr/local/bin/hotplug_monitor.sh"
